# roles/backup/files/restic-backup.timer
[Unit]
Description=Timer for daily backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
